# root/etc/services.d/squid/run
#!/usr/bin/execlineb -P

/usr/sbin/squid -NC
